# modules/ClusterConfig/playbooks/roles/nomad/files/nomad.conf
[program:nomad]
command=/usr/local/bin/nomad agent -config /etc/nomad.d
autostart=true
autorestart=true
stderr_logfile=/var/log/nomad.err.log
stdout_logfile=/var/log/nomad.out.log
